# packages/net/chrony/files/systemd/chronyd.service
[Unit]
Description=NTP client/server
After=ntpdate.service
Conflicts=ntpd.service

[Service]
ExecStart=/usr/sbin/chronyd -u ntp -d

[Install]
WantedBy=multi-user.target
