# packages/sys-devel/icecream/files/systemd/icecc-scheduler.service
[Unit]
Description=Icecream Distributed Compiler Scheduler
After=network.target

[Service]
User=icecc
EnvironmentFile=/etc/conf.d/icecream.conf
ExecStart=/usr/host/bin/icecc-scheduler $NETNAME -l /var/log/icecream/icecc-scheduler.log

[Install]
WantedBy=multi-user.target
