# packages/sys-apps/rng-tools/files/systemd/rngd.service
[Unit]
Description=Hardware RNG Entropy Gatherer Daemon

[Service]
ExecStart=@sbindir@/rngd -f

[Install]
WantedBy=multi-user.target
